# control_panel/templates/service-template.service
[Unit]
Description=Control Panel managed service: %i
After=network.target

[Service]
Type=simple
EnvironmentFile=%h/.config/control-panel/env/%i.env
ExecStart=/bin/bash -c "cd \"${WORKING_DIR}\" && ${COMMAND}"
Restart=on-failure

[Install]
WantedBy=default.target
